# systemd/hermes-profile.service.template
[Unit]
Description=Hermes {{PROFILE_NAME}} Discord gateway
After=network-online.target 9router.service
Wants=network-online.target

[Service]
Type=simple
User={{HERMES_USER}}
Group={{HERMES_USER}}
WorkingDirectory=/home/{{HERMES_USER}}/workspace
Environment=HERMES_PROFILE={{PROFILE_NAME}}
EnvironmentFile=/home/{{HERMES_USER}}/.hermes/profiles/{{PROFILE_NAME}}/.env
ExecStart=/home/{{HERMES_USER}}/.local/bin/hermes --profile {{PROFILE_NAME}} gateway start
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
